# debian/cyber-mariadb-zfs-snapshot.timer
[Timer]
OnCalendar=22:30

[Install]
WantedBy=timers.target
